# deploy/systemd/dn42-autopeer-agent.service
[Unit]
Description=dn42 Autopeer Node Agent
After=network-online.target bird.service
Wants=network-online.target

[Service]
# All settings (including the bearer token and WireGuard private key) live in the config file.
# Copy agent/config.example.json there and restrict it: chown root:root, chmod 0600.
ExecStart=/usr/local/bin/dn42-autopeer-agent -config /etc/dn42-autopeer/agent.json
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target
