# deploy/systemd/dn42-autopeer-bot.service
[Unit]
Description=dn42 Autopeer Telegram Bot
After=network-online.target dn42-autopeer-backend.service
Wants=network-online.target

[Service]
WorkingDirectory=/opt/dn42-autopeer/backend
EnvironmentFile=/opt/dn42-autopeer/backend/.env
ExecStart=/opt/dn42-autopeer/backend/.venv/bin/python -m app.bot.main
Restart=on-failure
User=autopeer
Group=autopeer

[Install]
WantedBy=multi-user.target
