# deploy/systemd/dn42-autopeer-backend.service
[Unit]
Description=dn42 Autopeer Backend
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/dn42-autopeer/backend
Environment=HOST=127.0.0.1
Environment=PORT=8000
EnvironmentFile=/opt/dn42-autopeer/backend/.env
ExecStart=/opt/dn42-autopeer/backend/.venv/bin/uvicorn app.main:app --host ${HOST} --port ${PORT}
Restart=on-failure
User=autopeer
Group=autopeer

[Install]
WantedBy=multi-user.target
